# scripts/check_wlan.service
[Unit]
Description=Check Wlan
 
[Service]
ExecStart=/usr/local/bin/check_wlan.sh
 
[Install]
WantedBy=default.target